# soft-hwclock-tick.service
[Unit]
Description=Save hwclock by timer

[Service]
Type=oneshot
ExecStart=/usr/bin/soft-hwclock tick > /dev/null
